# files/usr/lib/systemd/system/remove-fstab-root.service
[Unit]
Description=Remove root fstab entry as ostree already covers root
DefaultDependencies=no
Before=systemd-remount-fs.service
ConditionPathExists=/etc/fstab

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/remove-fstab-root

[Install]
WantedBy=systemd-remount-fs.service
